# builds/postgresql/container/postgres.ini
[program:postgresql]
command=/opt/init_postgres.sh
process_name=%(program_name)s
numprocs=1
directory=/var/lib/postgres
autostart=true
autorestart=true
startsecs=10
startretries=3
user=postgres
stdout_logfile=/var/log/supervisor/postgres-stdout.log
stderr_logfile=/var/log/supervisor/postgres-stderr.log